# etc/service/runsvdir/run
#!/bin/bash
/opt/opscode/embedded/bin/runsvdir-start
